# debian/libra-pool-mgm.upstart
description "Libra pool-mgm service"
author "Endre Karlson <endre.karlson@hp.com>"

start on runlevel [2345]
stop on runlevel [016]

respawn

pre-start script
    mkdir -p -m0755 /var/run/libra
    chown libra /var/run/libra
end script

exec /usr/bin/libra_pool_mgm -c /etc/libra.cfg --user=libra
