# debian/libra-admin-api.upstart
description "Libra admin-api service"
author "Endre Karlson <endre.karlson@hp.com>"

start on runlevel [2345]
stop on runlevel [016]

respawn

pre-start script
    mkdir -p -m0755 /var/run/libra
    chown libra /var/run/libra
end script

script
    [ -r "UPSTART_JOB" ] && {
        source $UPSTART_JOB
    }
    exec /usr/bin/libra_admin_api $LIBRA_OPTS
end script